# contrib/init/azlorad.conf
description "Azlora Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env AZLORAD_BIN="/usr/bin/azlorad"
env AZLORAD_USER="azlora"
env AZLORAD_GROUP="azlora"
env AZLORAD_PIDDIR="/var/run/azlorad"
# upstart can't handle variables constructed with other variables
env AZLORAD_PIDFILE="/var/run/azlorad/azlorad.pid"
env AZLORAD_CONFIGFILE="/etc/azlora/azlora.conf"
env AZLORAD_DATADIR="/var/lib/azlorad"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # azlorad will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$AZLORAD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run azlorad."
        echo "The setting must appear in $AZLORAD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Azlora Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$AZLORAD_PIDDIR"
    chmod 0755 "$AZLORAD_PIDDIR"
    chown $AZLORAD_USER:$AZLORAD_GROUP "$AZLORAD_PIDDIR"
    chown $AZLORAD_USER:$AZLORAD_GROUP "$AZLORAD_CONFIGFILE"
    chmod 0660 "$AZLORAD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$AZLORAD_PIDFILE" \
    --chuid $AZLORAD_USER:$AZLORAD_GROUP \
    --exec "$AZLORAD_BIN" \
    -- \
    -pid="$AZLORAD_PIDFILE" \
    -conf="$AZLORAD_CONFIGFILE" \
    -datadir="$AZLORAD_DATADIR" \
    -disablewallet \
    -daemon

